# roles/haproxy/files/haproxy-ssl.service.j2
# {{ ansible_managed }}
[Unit]
Description=Service for renewing haproxy certs
After=network-online.target

[Service]
ExecStart=/etc/systemd/haproxy-ssl.sh
Type=oneshot
